# .puppet/templates/ipfw.rules
KERNEL=="lo*", RUN+="/sbin/modprobe ipfw_mod"
